# linux-desktop/symlinks/keyd/dot-config/systemd/user/keyd-application-mapper.service
[Unit]
Description=Keyd Application Mapper

[Service]
ExecStart=keyd-application-mapper
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
